# conf/php/crontab.cron
* * * * * root cd /www/public_html/ ; php process.php
